# Buoi6/file_conf/runupload.conf
[program: uploaddata]
command=/home/pi/Desktop/buoi6/file_sh/upload.sh
autostart=true
autorestart=true
startsecs=60
stderr_logfile=/home/pi/Desktop/buoi6/log/up_error.log
stdout_logfile=/home/pi/Desktop/buoi6/log/up_logfile.log